# domestic-billing/deploy/iccorehub-domestic-billing.service
[Unit]
Description=ICCoreHub Domestic Billing Portal
After=network.target

[Service]
Type=simple
WorkingDirectory=/srv/iccorehub/domestic-billing
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=8788
ExecStart=/usr/bin/node /srv/iccorehub/domestic-billing/server.mjs
Restart=always
RestartSec=3
User=www-data
Group=www-data

[Install]
WantedBy=multi-user.target
